# systemd/bme280-collector.timer
[Unit]
Description=bme280 collector timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target
