# systemd/bme280-collector.service
[Unit]
Description=bme280 collector

[Service]
Type=oneshot
ExecStart=/usr/bin/bme280_collector
